# packages/vfuse-react/vfuse-dev-frontend-service.service
[Unit]
Description=VFuse Development Frontend Service
After=network.target
StartLimitIntervalSec=0
[Service]
Type=simple
Restart=always
RestartSec=1
User=isis
ExecStart=/bin/bash /home/isis/vfuse/js-vfuse/packages/vfuse-react PORT=80 HOST=0.0.0.0 npm run start

[Install]
WantedBy=multi-user.target

