# packages/vfuse-react/vfuse-dev-frontend.service
[Unit]
Description=VFuse Development Frontend Service
After=syslog.target network.target
StartLimitIntervalSec=0
[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/usr/bin/bash /home/isis/vfuse/js-vfuse/packages/vfuse-react/start.sh

[Install]
WantedBy=multi-user.target
