# packages/vfuse-react/vfuse-dev-frontend.service
[Unit]
Description=VFuse Development Frontend Service
After=network.target
StartLimitIntervalSec=0
[Service]
Type=simple
Restart=always
RestartSec=1
User=isis
ExecStart=PORT=80 HOST=0.0.0.0 /usr/bin/npm home/isis/vfuse/js-vfuse/packages/vfuse-react run start

[Install]
WantedBy=multi-user.target

