# packages/vfuse-bootstrap/vfuse-bootstrap.service
[Unit]
Description=VFuse Bootstrap Service
After=network.target
StartLimitIntervalSec=0
[Service]
Type=simple
Restart=always
RestartSec=1
User=isis
ExecStart=/bin/bash /usr/bin/node /home/isis/vfuse/js-vfuse/packages/vfuse-bootstrap/index.js

[Install]
WantedBy=multi-user.target
